# deployment/discord-ai-bot@.service
[Unit]
Description=Discord AI Agent Bot (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=discord-bot
Group=discord-bot
WorkingDirectory=/opt/discord-ai-bot
Environment="PATH=/opt/discord-ai-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
# ワークスペースのルートディレクトリを指定（オプション）
# Environment="AGENT_WORKSPACE_ROOT=/data/agent-workspaces"
EnvironmentFile=/opt/discord-ai-bot/.env
ExecStart=/opt/discord-ai-bot/venv/bin/python run.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
